# autostart/autotiling.desktop
[Desktop Entry]
Type=Application
Name=Autotiling
Exec=autotiling --limit 2

