# autostart/picom.desktop
[Desktop Entry]
Type=Application
Name=Picom
Exec=picom
