# systemd/user/caps-to-escape.service
[Unit]
Description=Set Caps Lock to Escape

[Service]
ExecStart=/usr/bin/setxkbmap -option "caps:escape"

[Install]
WantedBy=default.target
